# content/2016/2016-08-28-como-configurar-una-impresora-hp-en-red-en-arch-linux/code/install.sh
$ yaourt -S cups hplip
# systemctl start org.cups.cupsd.service
# systemctl start avahi-daemon.service
# systemctl enable org.cups.cupsd.service
# systemctl enable avahi-daemon.service