# misc/2013/gists/4573882/crontab
*/5 * * * * /home/raspberrypi/scripts/freedns_update.sh